# 91-mydlp-usb.rules
SUBSYSTEMS=="usb", KERNEL=="sd??", ACTION=="add", RUN+="/home/ozgur/usbtest.sh %k"
